# rise-temp6.service
[Unit]
Description=Rise Mode Temp 6 display
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/bin/rise_temp6.py watch --interval 1 --method feature
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
